# monitor.timer
[Unit]
Description=Run Deklan Node Monitor periodically
After=network-online.target
Wants=network-online.target

[Timer]
# Delay run after boot
OnBootSec=5m

# Default interval — overridden via install.sh
OnUnitActiveSec=3h

# Reduce API spike / system congestion
RandomizedDelaySec=45

# Maintain timer schedule across reboots
Persistent=true

# Target unit
Unit=monitor.service

[Install]
WantedBy=timers.target
