# monitor.timer
[Unit]
Description=Run Deklan Node Monitor every 3 hours
After=network-online.target
Wants=network-online.target

[Timer]
OnBootSec=5m
OnUnitActiveSec=3h
Persistent=true

[Install]
WantedBy=timers.target
