# monitor.service
[Unit]
Description=Deklan Node Monitor (oneshot)
After=network-online.target
Wants=network-online.target

# Pastikan file ada
ConditionPathExists=/opt/deklan-node-bot/monitor.py
ConditionPathExists=/opt/deklan-node-bot/.venv/bin/python

[Service]
Type=oneshot
User=root
Group=root

# ENV
EnvironmentFile=-/opt/deklan-node-bot/.env

# Working dir
WorkingDirectory=/opt/deklan-node-bot

# Smart Python runner
ExecStart=/bin/bash -c '
  if [ -x "/opt/deklan-node-bot/.venv/bin/python" ]; then
      exec /opt/deklan-node-bot/.venv/bin/python /opt/deklan-node-bot/monitor.py;
  else
      exec /usr/bin/python3 /opt/deklan-node-bot/monitor.py;
  fi
'

# Logging
StandardOutput=journal
StandardError=journal

Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=UTF-8"

# Timeout
TimeoutSec=90
KillMode=process

# Security (safe)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# oneshot → does not restart
Restart=no

[Install]
WantedBy=multi-user.target
