# monitor.timer
[Unit]
Description=Run Deklan Node Monitor periodically
After=network-online.target
Wants=network-online.target

[Timer]
# Delay after boot
OnBootSec=2m

# Default 3h
OnUnitActiveSec=3h

# Random spread (avoid many nodes hitting at same time)
RandomizedDelaySec=45

Persistent=true

Unit=monitor.service

[Install]
WantedBy=timers.target
