# monitor.timer
[Unit]
Description=Run Deklan Node Monitor periodically
After=network-online.target
Wants=network-online.target

[Timer]
# Delay execution after boot
OnBootSec=90

# Main interval — dynamically updated by install.sh based on .env:
#   MONITOR_EVERY_MINUTES
# Default fallback here:
OnUnitActiveSec=3h

# Anti–cluster spikes
RandomizedDelaySec=45

# Keep timer schedule across reboot
Persistent=true

# Target service
Unit=monitor.service

[Install]
WantedBy=timers.target
