# monitor.timer
[Unit]
Description=Run Deklan Node Monitor periodically
After=network-online.target
Wants=network-online.target

[Timer]
# Start 5 minutes after boot
OnBootSec=5m

# Default interval — will be updated by install.sh if env exists
OnUnitActiveSec=3h

# Spread monitor runs to reduce API burst
RandomizedDelaySec=45

# Remember schedule even after reboot
Persistent=true

Unit=monitor.service

[Install]
WantedBy=timers.target
